# conf/cron.d/update_plotly_graphs
*/3 * * * *    pi    cd /home/pi/waterino && ./make_graph_from_csv.py
